# systemd/nginx-config-reload.path
[Unit]
Description=Watch nginx whitelist config for changes

[Path]
PathModified=/etc/nginx/whitelist/*.conf

[Install]
WantedBy=multi-user.target